# bots/task34-bot/com.bridge34.task34bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bridge34.task34bot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/fireant/.openclaw/workspace/bots/task34-bot/.venv/bin/python</string>
    <string>/Users/fireant/.openclaw/workspace/bots/task34-bot/bot.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/fireant/.openclaw/workspace/bots/task34-bot</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/fireant/.openclaw/workspace/bots/task34-bot/task34bot.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/fireant/.openclaw/workspace/bots/task34-bot/task34bot.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>TZ</key>
    <string>Asia/Seoul</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
